# config-systemd-user/normalprofile.timer
[Unit]
Description=Turn on normal profile
[Timer]
Persistent=true
OnCalendar=*-*-* 6:30:00
[Install]
WantedBy=timers.target
