# config-systemd-user/normalprofile.service
[Unit]
Description=User Wake Up Action
[Service]
ExecStart=/home/adam/bin/profile-swap restore
