# config-systemd-user/silentprofile.timer
[Unit]
Description=Turn on silent profile
[Timer]
Persistent=true
OnCalendar=*-*-* 20:00:00
[Install]
WantedBy=timers.target
